# config/supervisord.d/mongod.ini
[program:mongod]
redirect_stderr = true
command = /usr/bin/mongod
process_name = mongod
