# config/supervisord.d/uwsgi.ini
[program:uwsgi]
command = /usr/sbin/uwsgi --ini /etc/uwsgi.ini
program_name = uwsgi
