# alerta-db/mongod.ini
[program:mongod]
redirect_stderr = true
command = /usr/bin/mongod -f /etc/mongod.conf
process_name = mongod
