# alerta-api/supervisord.d/uwsgi.ini
[program:uwsgi]
command = /usr/bin/uwsgi --ini /etc/uwsgi.ini
program_name = uwsgi
